# systemd-user/deck-before-sleep.service
[Unit]
Description=StreamDeck: set page 0 right before system sleep (login1 watcher)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
# Contexte X11 de ta session
Environment=DISPLAY=:1
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
ExecStart=%h/.local/bin/deck-sleep-watch
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
