# systemd-user/streamdeck-ui.service
[Unit]
Description=StreamDeck UI (streamdeck-linux-gui)
# Démarrer après l'ouverture de session graphique
After=graphical-session.target
PartOf=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple

# Contexte graphique X11 (tes valeurs actuelles)
Environment=QT_QPA_PLATFORM=xcb
Environment=DISPLAY=:1
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority

# Petit délai pour laisser la session graphique se stabiliser
ExecStartPre=/usr/bin/sleep 1

# Chemin du binaire (pipx installe un wrapper ici)
ExecStart=%h/.local/bin/streamdeck -n

Restart=always
RestartSec=2

[Install]
# Lie l’unité à la session graphique (pas à default.target)
WantedBy=graphical-session.target
